# theWeather_system.service
[UNIT]
Description=Start theWeather.system statistics

[Service]
ExecStart=/usr/bin/startProcedure.sh
KillMode=process

[Install]
WantedBy=multi-user.target